# rootfs/usr/lib/systemd/user/sunshine.service.d/99-realtime.conf
[Service]
# Real-Time First-In-First-Out Scheduling
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=10

# High process priority
Nice=-15

# Real-Time I/O (Disk access)
IOSchedulingClass=realtime
IOSchedulingPriority=0
