# rootfs/usr/lib/systemd/user/sunshine.service.d/99-realtime.conf
[Service]
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=10
Nice=-15
IOSchedulingClass=realtime
IOSchedulingPriority=0
